# test/stubs/crontab
#!/bin/bash
echo "crontab $*" >> "${STUB_LOG:-/tmp/claude-catcher-stub.log}"
